# utilities/trv-pi-001/services/journal-to-tsstore.service
[Unit]
Description=Stream journalctl to tsstore
After=tsstore.service

[Service]
Type=simple
ExecStart=/home/tviviano/tsstore/scripts/journal-to-tsstore.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
